# services/sillytavern.service
[Unit]
Description=SillyTavern Chat Interface
After=network.target llama-server.service
RequiresMounts=/media/ryan/UbuntuDATA

[Service]
Type=simple
ExecStart=/bin/bash /media/ryan/UbuntuDATA/bin/sillytavern-wrapper.sh
WorkingDirectory=/media/ryan/UbuntuDATA/AI_PROJECTS/SillyTavern-home
Restart=on-failure
RestartSec=5
TimeoutStopSec=10
KillMode=mixed
KillSignal=SIGTERM
SendSIGKILL=yes

[Install]
WantedBy=default.target
